# apps/pi-forwarder/systemd/loramapr-pi-forwarder.service
[Unit]
Description=LoraMapr Pi Forwarder
After=network-online.target
Wants=network-online.target

[Service]
Type=simple
User=pi
Group=pi
WorkingDirectory=/opt/loramapr/pi-forwarder
EnvironmentFile=/etc/loramapr/pi-forwarder.env
ExecStart=/usr/bin/node /opt/loramapr/pi-forwarder/dist/index.js
Restart=always
RestartSec=3
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
